# Pinbot/deploy_conf/supervisor/supervisord.conf
; Sample supervisor config file.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Note: shell expansion ("~" or "$HOME") is not supported.  Environment
; variables can be expanded using this syntax: "%(ENV_HOME)s".

[unix_http_server]
file=/data/pinbot/tmp/pinbot_supervisor.sock   ; (the path to the socket file)
chmod=0700                 ; socket file mode (default 0700)
chown=deploy:deploy       ; socket file uid:gid owner

[inet_http_server]         ; inet (TCP) server disabled by default
port=127.0.0.1:9004        ; (ip_address:port specifier, *:port for all iface)
username=pinbot              ; (default is no username (open server))
password=hopperclouds2013               ; (default is no password (open server))

[supervisord]
logfile=/data/pinbot/log/supervisord/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
pidfile=/data/pinbot/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false               ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
;umask=022                   ; (process file creation umask;default 022)
;user=chrism                 ; (default is current user, required if root)
;identifier=supervisor       ; (supervisord identifier, default is 'supervisor')
;directory=/tmp              ; (default is not to cd during start)
;nocleanup=true              ; (don't clean up tempfiles at start;default false)
;childlogdir=/tmp            ; ('AUTO' child log dir, default $TEMP)
;environment=KEY=value       ; (key value pairs to add to environment)
;strip_ansi=false            ; (strip ansi escape codes in logs; def. false)

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///data/pinbot/tmp/pinbot_supervisor.sock ; use a unix:// URL  for a unix socket

[program:web]
command=/home/deploy/Pinbot/pin_venv/bin/uwsgi /home/deploy/Pinbot/deploy_conf/uwsgi/uwsgi.ini
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/home/deploy/Pinbot                ; directory to cwd to before exec (def no cwd)
autorestart=true        ; whether/when to restart (default: unexpected)
autorestart=true
stdout_logfile=/data/pinbot/log/supervisord/pinbot_web.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
redirect_stderr=true
stderr_logfile=/data/pinbot/log/supervisord/pinbot_web_err.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
stopsignal=HUP

[program:c_worker]
environment=PYTHONPATH=/home/deploy/Pinbot
command=/home/deploy/Pinbot/pin_venv/bin/python pin_celery/celery_app.py worker -l info -n web_worker --concurrency=8
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/home/deploy/Pinbot                ; directory to cwd to before exec (def no cwd)
autorestart=unexpected        ; whether/when to restart (default: unexpected)
stdout_logfile=/data/pinbot/log/supervisord/celery_worker.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
redirect_stderr=true
stderr_logfile=/data/pinbot/log/supervisord/celery_err.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)

[program:c_beat]
environment=PYTHONPATH=/home/deploy/Pinbot
command=python pin_celery/celery_app.py beat -s /data/pinbot/log/supervisord/celerybeat-schedule --pidfile=/data/pinbot/log/supervisord/beat.pid
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/home/deploy/Pinbot                ; directory to cwd to before exec (def no cwd)
autorestart=unexpected        ; whether/when to restart (default: unexpected)
stdout_logfile=/data/pinbot/log/supervisord/celery_beat.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
redirect_stderr=true
stderr_logfile=/data/pinbot/log/supervisord/celery_beat_err.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)

[program:flower]
environment=PYTHONPATH=/home/deploy/Pinbot
command=pin_venv/bin/celery flower --basic_auth=admin:Hopperclouds2014 --broker=amqp://admin:Rabbitmq2015@localhost:5672// --port=5555
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
directory=/home/deploy/Pinbot                ; directory to cwd to before exec (def no cwd)
autorestart=unexpected        ; whether/when to restart (default: unexpected)
stdout_logfile=/data/pinbot/log/supervisor/flower.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stderr_logfile=/data/pinbot/log/supervisor/flower_err.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=50MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)
stopsignal=QUIT
